# deployment/systemd/rad.io.trigger.path
[Unit]
Description=Watch for RAD.io instance

[Path]
PathChanged=<project repo path>/src/
PathChanged=<project repo path>/bun.lock

[Install]
WantedBy=multi-user.target
